# scripts/pixconvert.service
[Unit]
Description=PixConvert API Server
After=network.target

[Service]
Type=simple
User=www-data
ExecStart=/usr/bin/node /usr/local/lib/node_modules/pixconvert/server.js
ExecStartPre=/bin/mkdir -p /var/log/pixconvert
Restart=on-failure
RestartSec=5

# Environment — edit as needed
Environment=PORT=3001
Environment=API_RATE_LIMIT=10
Environment=FILE_SIZE_LIMIT_MB=50
Environment=FILE_TTL_HOURS=1

# Logging
StandardOutput=append:/var/log/pixconvert/pixconvert_access.log
StandardError=append:/var/log/pixconvert/pixconvert_error.log
SyslogIdentifier=pixconvert

[Install]
WantedBy=multi-user.target
